# one/supervisord.conf
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:k8swatcher]
command=python -u k8s_watcher.py
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/k8swatcher_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/k8swatcher_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",LAUNCHER_COUNT="1",MODULE_NAME="k8swatcher"

[program:launcher]
command=python -u launcher.py
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/launcher_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/launcher_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",LAUNCHER_COUNT="1",MODULE_NAME="launcher"

[program:scheduler]
command=python -u scheduler.py
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/scheduler_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/scheduler_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",LAUNCHER_COUNT="1",MODULE_NAME="scheduler"

[program:query_server]
command=python -u uvicorn_server.py --port 8081
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/query_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/query_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",MODULE_NAME="query-server",SERVER="query"

[program:operating_server]
command=python -u uvicorn_server.py --port 8082
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/operating_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/operating_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",MODULE_NAME="server",SERVER="operating"

[program:ugc_server]
command=python -u uvicorn_server.py --port 8083
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/ugc_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/ugc_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",MODULE_NAME="ugc-server",SERVER="ugc"

[program:monitor_server]
command=python -u uvicorn_server.py --port 8084
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/high-flyer/log/monitor_0.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/monitor_0.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/high-flyer/code/multi_gpu_runner_server",MARSV2_MANAGER_CONFIG_DIR="/etc/hai_one_config",REPLICA_RANK="0",MODULE_NAME="monitor-server",SERVER="monitor"

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg
directory=/high-flyer/code/multi_gpu_runner_server
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
redirect_stderr=true
stdout_logfile=/high-flyer/log/haproxy.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/haproxy.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

[program:studio]
command=/marsv2/scripts/studio/hai-studio --log-dir /high-flyer/log
priority=999
startsecs=1
startretries=6
autostart=true
autorestart=unexpected
exitcodes=0
redirect_stderr=true
stdout_logfile=/high-flyer/log/studio.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/high-flyer/log/studio.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
